# deploy/sidecar.service
[Unit]
Description=Sidecar DMS - Unified backend (HTTP + intake watcher)
After=network.target

[Service]
Type=simple
User=jodli
WorkingDirectory=/mnt/quickstuff/git/sidecar-dms
EnvironmentFile=/mnt/quickstuff/git/sidecar-dms/.env
ExecStart=/mnt/quickstuff/git/sidecar-dms/.venv/bin/python tools/server.py
Restart=on-failure
RestartSec=10

# Graceful shutdown: uvicorn propagates SIGTERM, watcher task stops cleanly
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
